# env/puppet/modules/buddyup/templates/gunicorn.service.erb
[Unit]
Description=Starts the gunicorn server for <%= name %>

[Service]
ExecStart=<%= daemon_command %>